# ping-plotter.conf
description "ping plotter pinger"

start on runlevel [23]
stop on shutdown

script
  cd /apps/ping-plotter
  exec sudo /usr/local/bin/ruby ping.rb
end script

respawn
