# ping-plotter-server.conf
description "ping plotter server"

start on runlevel [23]
stop on shutdown

script
  cd /apps/ping-plotter
  exec sudo thin start -p 9988
end script

respawn
